# templates/pulp_celerybeat.service.erb
[Unit]
Description=Pulp Celerybeat
Wants=syslog.service

[Service]
Restart=always
RestartSec=30
TimeoutStartSec=0
TimeoutSec=300
ExecStartPre=-/usr/bin/podman pull quay.io/foreman/pulp-celerybeat:latest
ExecStartPre=-/usr/bin/podman rm "pulp-celerybeat-1"
ExecStart=/usr/bin/podman run --name pulp-celerybeat-1 --net host -v /var/lib/pulp:/var/lib/pulp:Z -v /etc/pulp:/etc/pulp:Z quay.io/foreman/pulp-celerybeat:latest
ExecReload=-/usr/bin/podman stop "pulp-celerybeat-1"
ExecReload=-/usr/bin/podman rm "pulp-celerybeat-1"
ExecStop=-/usr/bin/podman stop "pulp-celerybeat-1"
EnvironmentFile=-/etc/sysconfig/pulp-celerybeat

[Install]
WantedBy=multi-user.target
